# store_app.conf
description "uWSGI server instance configured to server store_app"

start on runlevel [2345]
stop on runlevel [!2345]

setuid vagrant
setgid www-data

env PATH=/home/vagrant/backend/backend-env/bin
chdir /home/vagrant/backend/store_backend
exec uwsgi --ini store_app.ini